# mail/getmail.cron
# minute | hour | day_of_month | month | day_of_week | command

MAILTO=''
SHELL=/bin/bash
*/2 * * * 1-5 "$REPOS_BASE"/config/mail/getm
